# sh/supervisor/laravel-worker.conf
[program:portfolio-dem-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/html/portfolio-dem/artisan queue:work --sleep=1 --tries=3 --timeout=120 --queue=default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/portfolio-dem/storage/logs/worker.log
stopwaitsecs=3600
